# dispatch-router/ebench_sender.service
[Unit]
Description=Ebench sender service
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/root/start_senders.sh
