# dispatch-router/ebench_receiver.service
[Unit]
Description=Ebench receiver service
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/root/start_receivers.sh
